# deploy/amazon-ec2/circle-tasks.conf
#!upstart
#
# Copy this to /etc/init/circle-tasks.conf
# Modified from https://www.exratione.com
#
description "Circle Tasks project management server"

start on startup
stop on shutdown
respawn

# This line is needed so that Upstart reports the pid of the Node.js process
# started by Forever rather than Forever's pid.
expect fork

env NODE_BIN_DIR="/usr/local/bin"
env NODE_PATH="/usr/local/lib/node_modules"
env APPLICATION_DIRECTORY="/home/ubuntu/apps/circle-tasks/server"
env APPLICATION_START="server.js"
env LOG="/var/log/circle-tasks.log"
env AS_USER="start-stop-daemon --start -c ubuntu --exec"

pre-start script
    # Redirect port 80 requests to port 3000
    iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 80 -j REDIRECT --to-port 3000

    # Redirect port 443 requests to port 4000
    iptables -t nat -A PREROUTING -i eth0 -p tcp --dport 443 -j REDIRECT --to-port 4000

    # Give us permission to access the logs
    [ ! -f $LOG ] && { touch $LOG; }
    chown ubuntu:ubuntu $LOG
end script

script
    # Add the node executables to the path, which includes Forever if it is
    # installed globally, which it should be.
    PATH=$NODE_BIN_DIR:$PATH
    # The minUptime and spinSleepTime settings stop Forever from thrashing if
    # the application fails immediately on launch. This is generally necessary to
    # avoid loading development servers to the point of failure every time
    # someone makes an error in application initialization code, or bringing down
    # production servers the same way if a database or other critical service
    # suddenly becomes inaccessible.
    #
    # The w and watchDirectory options allows us to restart the app after
    # there have been changes.
    exec $AS_USER /usr/local/bin/forever -- --sourceDir $APPLICATION_DIRECTORY -a -l $LOG \
         --minUptime 5000 --spinSleepTime 2000 \
         -w --watchDirectory $APPLICATION_DIRECTORY start $APPLICATION_START
end script

pre-stop script
    # Add the node executables to the path.
    PATH=$NODE_BIN_DIR:$PATH
    # Here we're using the pre-stop script to stop the Node.js application
    # process so that Forever is given a chance to do its thing and tidy up
    # its data. Note that doing it this way means that each application that
    # runs under Forever must have a different start file name, regardless of
    # which directory it is in.
    exec $AS_USER /usr/local/bin/forever -- stop $APPLICATION_START >> $LOG
end script

post-stop script
    # Revert iptables rules
    iptables -t nat -D PREROUTING -i eth0 -p tcp --dport 80 -j REDIRECT --to-port 3000
    iptables -t nat -D PREROUTING -i eth0 -p tcp --dport 443 -j REDIRECT --to-port 4000
end script